# budgetin-bot.service
[Unit]
Description=Budgetin Telegram Bot
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/budgetin-bot
Environment=PATH=/home/ubuntu/budgetin-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/ubuntu/budgetin-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=budgetin-bot

[Install]
WantedBy=multi-user.target
